# docker/supervisor/supervisord.conf
; =============================================================================
; supervisord.conf — Gestión de Procesos en Producción
; Laravel 12 Docker Ready-to-Deploy Starter Kit
;
; OPTIMIZACIONES CLAVE:
;   - Un proceso supervisor único gestiona tanto PHP-FPM como los workers de
;     cola y el scheduler, evitando múltiples entrypoints en el contenedor.
;   - Los workers de cola usan '--tries' y '--timeout' para evitar jobs zombies
;     que bloqueen el proceso indefinidamente.
;   - 'stopasgroup=true' y 'killasgroup=true' garantizan que al detener un
;     worker, todos sus procesos hijo se terminan limpiamente (signal propagation).
;   - La rotación de logs está controlada por 'stdout_logfile_maxbytes' para
;     evitar que los logs llenen el filesystem del contenedor.
;   - El scheduler usa un sleep de 60 s en un bucle infinito, que es el patrón
;     recomendado para ejecutar 'artisan schedule:run' dentro de contenedores
;     (no hay cron daemon disponible por defecto en imágenes Alpine/Debian slim).
; =============================================================================

[unix_http_server]
; Socket para que 'supervisorctl' se comunique con el daemon
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
; --- Configuración global del daemon ---
nodaemon=true                       ; Corre en foreground para que Docker capture stdout/stderr
logfile=/dev/null                   ; Los logs del daemon propio van a stdout del contenedor
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid
user=root                           ; Supervisor necesita root para iniciar procesos como www-data

[rpcinterface:supervisor]
; Interfaz RPC requerida por supervisorctl
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; =============================================================================
; PROGRAMA: php-fpm
; Servidor FastCGI de PHP — núcleo de la aplicación Laravel
; =============================================================================
[program:php-fpm]
command=/usr/local/sbin/php-fpm --nodaemonize --force-stderr
; '--nodaemonize' es CRÍTICO: evita que PHP-FPM se demonize y escape
; al control de Supervisor. Sin esto, Supervisor lo marcará como "muerto".
autostart=true
autorestart=true
priority=10                         ; Arranca antes que los workers (número menor = mayor prioridad)
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=www-data

; =============================================================================
; PROGRAMA: laravel-worker
; Cola de trabajos — gestiona jobs asíncronos (emails, notificaciones, etc.)
;
; OPTIMIZACIONES CLAVE:
;   - 'numprocs=2': Dos workers en paralelo para throughput básico de producción.
;     Escalar según el volumen de jobs. En Kubernetes, preferir más réplicas de pod.
;   - '--sleep=3': Si no hay jobs, el worker duerme 3 s antes de volver a consultar.
;     Reduce carga en Redis/DB sin sacrificar latencia significativa.
;   - '--tries=3': Un job fallará definitivamente tras 3 intentos, evitando loops infinitos.
;   - '--timeout=60': Mata jobs que superan 60 s de ejecución. Previene workers zombies.
;   - '--max-jobs=1000': El worker se reinicia tras 1000 jobs para liberar memoria
;     acumulada (memory leaks de larga duración en PHP). Supervisor lo relanza limpio.
;   - '--max-time=3600': Reinicio forzado cada hora como medida de higiene adicional.
; =============================================================================
[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work %(ENV_QUEUE_CONNECTION)s \
    --queue=%(ENV_QUEUE_NAME)s \
    --sleep=3 \
    --tries=3 \
    --timeout=60 \
    --max-jobs=1000 \
    --max-time=3600
; El uso de variables de entorno (%(ENV_*)s) permite configurar la conexión
; y el nombre de la cola desde el .env sin tocar este archivo.
autostart=true
autorestart=true
numprocs=2
priority=20
startretries=5
startsecs=3                         ; El proceso debe vivir al menos 3 s para considerarse "iniciado"
stopwaitsecs=120                    ; Espera hasta 120 s para que el worker termine el job actual
; CRÍTICO para PHP: propagar señales a procesos hijo para evitar jobs cortados a mitad
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=www-data

; =============================================================================
; PROGRAMA: laravel-scheduler
; Task Scheduler de Laravel — equivalente al cron de producción
;
; OPTIMIZACIONES CLAVE:
;   - En contenedores, no existe un daemon cron. El patrón correcto es un
;     bucle infinito con 'sleep 60' que ejecuta 'schedule:run' cada minuto.
;   - Se usa 'bash -c' con un bucle 'while true' para que Supervisor gestione
;     el proceso padre y pueda terminarlo limpiamente.
;   - 'autorestart=unexpected' — solo reinicia si el proceso muere inesperadamente,
;     no en cada ciclo normal del sleep.
; =============================================================================
[program:laravel-scheduler]
command=bash -c "while true; do php /var/www/html/artisan schedule:run --no-interaction 2>&1; sleep 60; done"
autostart=true
autorestart=unexpected
priority=30
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=www-data
